# piccolo/0.0.1/sysext/piccolo/usr/lib/systemd/system/piccolo.service
[Unit]
Description=Symphony Piccolo Agent

[Service]
ExecStart=/usr/bin/piccolo
Restart=always

[Install]
WantedBy=multi-user.target
